# systemd/transit-tv-start-docker.service
[Unit]
Description=Transit TV Graceful Start (Docker)
Documentation=https://github.com/jasonad123/Transit-TV
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
# Use docker exec to run Node.js inside the container where localhost is actually localhost
# This bypasses the API's localhost-only security check and works reliably with Alpine
ExecStart=/usr/bin/docker exec transit-tv node -e "require('http').request({host:'localhost',port:8080,method:'POST',path:'/api/server/start'},r=>{r.on('data',d=>process.stdout.write(d));r.on('end',()=>process.exit(0))}).on('error',e=>{console.error(e);process.exit(1)}).end()"
RemainAfterExit=no
StandardOutput=journal
StandardError=journal
SyslogIdentifier=transit-tv-start

# Timeout settings
TimeoutStartSec=30s

[Install]
WantedBy=multi-user.target
